# etc/cron.daily/api_update.daily
#!/bin/sh

test -x /root/examples_update_silent.sh || exit 0
/root/examples_update_silent.sh
